# orchestrator/crontab
30 3 * * * musiclib --config /app/config.yaml --apply
